# deploy/us-stock-alert-bot.service.example
[Unit]
Description=US Stock Alert Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/usmart-quant
EnvironmentFile=/opt/usmart-quant/.env/finnhub.env
ExecStart=/opt/usmart-quant/.venv/bin/python -m stock_alert_bot.app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
